# ramdisk/rfs/etc/inittab
console::sysinit:/etc/init.d/rcS
#ttyS0::askfirst:-/bin/sh
::respawn:-/bin/HHlogin
::ctrlaltdel:/bin/umount -a -r
